# init/compiler-explorer.service
[Unit]
Description=Compiler Explorer

[Service]
Type=simple
Environment=HOME=/root
TimeoutStartSec=infinity
ExecStart=/compiler-explorer-image/init/start.sh

[Install]
WantedBy=multi-user.target
