# init/compiler-explorer.service
[Unit]
Description=Compiler Explorer
After=remote-fs.target systemd-logind.service network-online.target cloud-config.service
Wants=network-online.target
Requires=remote-fs.target

[Service]
Type=simple
User=root
TimeoutStartSec=infinity
WorkingDirectory=/infra
ExecStart=/infra/init/start.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=compiler-explorer

[Install]
WantedBy=multi-user.target
